# Sber - Data_Engineer/Module_Prof/final_project/main.cron
45 9 * * * /home/de11tm/ykir/project/main.py